# ops/systemd/parameter-golf-pr-enrichment.timer
[Unit]
Description=Run Parameter Golf data refresh every 30 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=30min
Unit=parameter-golf-pr-enrichment.service

[Install]
WantedBy=timers.target
